# files/afterInstall.service
[Unit]
Description=MeuScript

[Service]
ExecStart=/usr/bin/postInstall

[Install]
WantedBy=multi-user.target 